# debian/authd-example-broker.postinst
#!/bin/sh
set -e

#DEBHELPER#

# Mask the service by default, unless it was previuously enabled.
if [ -z "${DPKG_ROOT:-}" ] && [ -d /run/systemd/system ] && [ "$1" = configure ]; then
    if [ -z "$2" ] || ! deb-systemd-helper was‐enabled; then
        deb-systemd-helper mask 'authd-example-broker.service'
    fi
fi
